# ec2/crontab
*/1 * * * * ~/lichess-swiss-league-tournament-creation/script.sh start